# etc/template/api.ini
[program:recipe_api__ENV__]
user = __USERNAME__
directory = __APP_ROOT__
command = echo "Hello, world"

autorestart = true
redirect_stderr = true
stdout_logfile = __APP_ROOT__/var/log/%(program_name)s.log
stderr_logfile = __APP_ROOT__/var/log/%(program_name)s_err.log
