# etc/template/api.ini
[program:recipe_api_{{app_env}}]
user = {{username}}
directory = {{app_root}}
command = echo "Hello, world"

autorestart = true
redirect_stderr = true
stdout_logfile = {{app_root}}/var/log/%(program_name)s.log
stderr_logfile = {{app_root}}/var/log/%(program_name)s_err.log
